# develop-gdc-res-base-unit-tests-el9/mariadb/mariadb-finalize.service
[Unit]
After=mariadb.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'mysql < /opt/setup-unit.sql'

[Install]
WantedBy=multi-user.target
